# automation/launchd/com.genuverity.daily-research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.genuverity.daily-research</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/klop/GenuVerity7/automation/scripts/daily-auto-research.sh</string>
    </array>
    
    <!-- Run daily at 5:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/klop</string>
    </dict>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/klop/GenuVerity7/automation/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/klop/GenuVerity7/automation/logs/launchd-stderr.log</string>
    
    <!-- Run even if user not logged in (if possible) -->
    <key>RunAtLoad</key>
    <false/>
    
    <!-- Keep trying if it fails -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/klop/GenuVerity7/automation</string>
</dict>
</plist>
